# sparse/usr/lib/systemd/system/sailfish-device-encryption-community-wizard.service.d/20-wait-hwc.conf
[Service]
ExecStartPre=/usr/bin/droid/waitforservice init.svc.vendor.hwcomposer-2-3 running
